# OPTIMIZED-PYTHONANYWHERE/systemd/fight-agent-optimized.service
[Unit]
Description=Fight Agent Optimized Webhook Service
After=network.target

[Service]
Type=simple
User=trinity
Group=trinity
WorkingDirectory=/home/trinity/the-fight-predictor-agent/OPTIMIZED-PYTHONANYWHERE
EnvironmentFile=/home/trinity/the-fight-predictor-agent/OPTIMIZED-PYTHONANYWHERE/.env
ExecStart=/home/trinity/the-fight-predictor-agent/OPTIMIZED-PYTHONANYWHERE/.venv/bin/gunicorn --workers 1 --bind 127.0.0.1:8080 app:create_app()
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
